# overlay/etc/sv/nsca/run
#!/bin/bash

if [ -f "/opt/nsca/nsca.cfg" ]; then
    exec /usr/sbin/nsca -c /opt/nsca/nsca.cfg -f
fi
