# overlay/etc/sv/nsca/run
#!/bin/bash

if [ -f "/opt/nagios/etc/nsca.cfg" ]; then 
    exec /opt/nagios/bin/nsca -f -c /opt/nagios/etc/nsca.cfg --single
else
    echo "nsca configuration not found; skipping."
fi