# overlay/etc/sv/xinetd/run
#!/bin/bash

exec /usr/sbin/xinetd -f /etc/xinetd.conf -dontfork